# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

; ADK backend — internal only (127.0.0.1:8000)
; AWS credentials come from the environment (env vars or IAM role).
[program:adk]
command=adk web --host 127.0.0.1 --port 8000 .
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

; Next.js frontend — public (0.0.0.0:3000)
; HOSTNAME and PORT are set via Dockerfile ENV.
[program:nextjs]
command=node server.js
directory=/app/frontend
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
